# init.d/mopidy
#!/sbin/openrc-run

command="/media/mmcblk0p2/.mopidy/.local/bin/mopidy"
command_user="mopidy"
command_background="yes"
pidfile="/run/$RC_SVCNAME.pid"
logfile="/media/mmcblk0p2/logs/mopidy.log"
errfile="/media/mmcblk0p2/logs/mopidy.err"

depend() {
    need net
}

start() {
        rm -f "${logfile}" "${errfile}"
	ebegin "Starting mopidy"
	start-stop-daemon $daemon_args --start \
		--exec "${command}" \
                --user "${command_user}" \
                --make-pidfile \
		--pidfile "${pidfile}" \
                --background \
                --stdout "${logfile}" \
                --stderr "${errfile}" \
		-- ${ARGS}
	eend $? "Failed to start mopidy"
}

stop() {
	ebegin "Stopping mopidy"
	start-stop-daemon --stop --quiet \
		--pidfile "${pidfile}"
	eend $? "Failed to stop mopidy"
}

